# pi/system/camera-proxy.service
sudo tee /etc/systemd/system/camera-proxy.service << 'EOF'
[Unit]
Description=Camera Proxy Service
After=network.target

[Service]
ExecStart=/home/pi/venv/bin/python3 /home/pi/goatdev/pi/servers/camera_proxy.py
WorkingDirectory=/home/pi/goatdev/pi
User=root
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
EOF

sudo systemctl daemon-reload
sudo systemctl restart camera-proxy
systemctl status camera-proxy