# pi/system/camera-proxy.service
[Unit]
Description=Camera Proxy Service
After=network.target

[Service]
ExecStart=/home/pi/venv/bin/gunicorn -w 1 -k gevent --bind 0.0.0.0:8080 --timeout 300 servers.camera_proxy:app
WorkingDirectory=/home/pi/goatdev/pi
User=root
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
EOF

sudo systemctl daemon-reload
sudo systemctl restart camera-proxy
systemctl status camera-proxy