# pi/system/goat-heartbeat.cron
*/2 * * * * /home/pi/venv/bin/python /home/pi/goat-capture/pi/logger/pi_heartbeat_cron.py